# rootfs/etc/services.d/rakomqtt/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the rakomqtt service
# ==============================================================================

bashio::log.info "Starting Rako MQTT Bridge..."
exec /usr/src/app/start.sh
